# debian/postinstall.sh
cd ./lib/tprinter/

npm install --production

systemctl enable tprinter.service